# contest-bot.service
[Unit]
Description=Contest Reminder Telegram Bot
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/bot/directory
Environment="BOT_TOKEN=your_bot_token_here"
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
